# codedeploy-zero-downtime-systemd/golang/deploy-tools/systemd/system/web.service
[Service]
Type=simple
Environment="STARTED_FROM_SYSTEMD=1"
ExecStart=/home/ec2-user/bin/application
