# elastic-beanstalk-zero-downtime-systemd/golang/deploy-tools/systemd/system/web_2.service
[Service]
Type=simple
Environment="PORT=5000"
Environment="STARTED_FROM_SYSTEMD=1"
ExecStart=/var/app/current/bin/application
